# home/update.cron
@reboot /home/roger/update.sh
0 4 * * 1 /home/roger/update.sh
